# moztrap-supervisor.conf

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"

[program:uwsgi]
command = /usr/local/bin/uwsgi /moztrap-uwsgi.ini

